# chroma-addon/rootfs/etc/services.d/chroma/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# Take down the S6 supervision tree when ChromaDB fails
# ==============================================================================

# First argument is crashcode, second is signal
# If crashcode is not 0, and signal is not SIGTERM (15), then something went wrong
if-then
    # Only if both conditions are true, we should drop S6
    elglob -0 -- crashcode { 0 }
    elglob -0 -- signal { 15 }
    foreground { s6-svscanctl -t /var/run/s6/services }
    if { s6-echo "[finish] ChromaDB exited with code ${crashcode} and signal ${signal}" }

foreground { s6-svscanctl -t /var/run/s6/services }
